# discord-bot.service
[Unit]
Description=Discord Bot - Hybrid AI System
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=jarvis
Group=bots
WorkingDirectory=/opt/jarvis/discord-bot
Environment=PATH=/opt/jarvis/discord-bot/venv/bin
ExecStart=/opt/jarvis/start.sh
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bot

# Security settings (relaxed for compatibility)
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=1024
MemoryMax=512M

[Install]
WantedBy=multi-user.target